# scripts/powervpn.service
[Unit]
Description=PowerVPN Management Panel
After=network.target mysql.service

[Service]
# User and Group
User=vpnpanel
Group=vpnpanel
WorkingDirectory=/opt/powervpn

# Execution
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=5

# Defect 47: Hardening
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
ReadOnlyPaths=/
ReadWritePaths=/opt/powervpn/.next /opt/powervpn/logs
CapabilityBoundingSet=
RestrictRealtime=yes
ProtectControlGroups=yes
ProtectKernelModules=yes

# Limits
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
